# docker/s6-supervisor/services.d/php-fpm/run
#!/usr/bin/env bash

echo "PHP-FPM is starting..."
exec /usr/local/sbin/php-fpm --nodaemonize